# 11-Scripts/proxmox/mount-monitor/proxmox-mount-monitor.service
[Unit]
Description=Mount Watchdog (LXC + VM)
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/proxmox-mount-monitor.sh
