# 11-Scripts/proxmox/smb-vm-mount.service
[Unit]
Description=Mount CIFS share from VM 200 after VM starts
Wants=network-online.target // wants network to be available
After=network-online.target // starts only after network is up

[Service]
Type=oneshot // run once and exit
ExecStartPre=/bin/bash -c 'until ping -c1 192.168.2.200; do sleep 2; done' // wait until VM is reachable
ExecStart=/root/mount-torrent.sh // execute the mount script
RemainAfterExit=yes
TimeoutSec=0

[Install]
WantedBy=multi-user.target // start automatically at multi-user runlevel
